# hello.service
[Unit]
Description=Hello Timer Service

[Service]
Type=oneshot
ExecStart=/usr/bin/bash -c 'echo "Hello from timer at $(date)" >> /var/log/hello.log'
