# hello.service
[Unit]
Description=Hello Timer Service

[Service]
Type=oneshot
ExecStart=/usr/bin/bash -c '/app/hello > /var/log/hello.log'
